# deploy/arcana-world.service
[Unit]
Description=arcana world engine — capability-fenced runtime WASM generation (gen-server)
Documentation=https://github.com/jrjohn/wasm-jit
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=arcana
Group=arcana
WorkingDirectory=/opt/arcana/app
ExecStart=/opt/arcana/app/gen-server
# GOOGLE_CLIENT_ID and CLAUDE_CODE_OAUTH_TOKEN live here, mode 0600, root-owned.
# Without GOOGLE_CLIENT_ID the server runs in local mode and every write door is
# open — on a public address that would be the whole point, undone. Setting it is
# what turns the gate on.
EnvironmentFile=/etc/arcana/world.env
Restart=on-failure
RestartSec=3

# The engine binds loopback only; nginx is the sole way in from outside.
# Confinement here is the same idea as the cell's: hand it exactly what it needs.
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
LockPersonality=yes
# Everything visitors create — worlds, the ālaya ledger, grown skins and widgets.
ReadWritePaths=/opt/arcana/app/worlds /opt/arcana/app/ledger /opt/arcana/app/skins-grown /opt/arcana/app/widgets-grown

[Install]
WantedBy=multi-user.target
